# crontab
*/1 * * * * python /usr/src/app/index.py -a -v -o /data
